# build/supervisor/promtail.conf
[group:promtail]
programs=promtail-server

[program:promtail-server]
directory=/app
command=/usr/local/bin/promtail --config.file=/etc/promtail/config-promtail.yml -config.expand-env=true
user=root
autostart=true
autorestart=true
stopsignal=TERM

stderr_logfile_maxbytes=1024MB
stdout_logfile_maxbytes=1024MB
stdout_logfile_backups=2
stdout_logfile=/data/logs/promtail.log
stderr_logfile=/data/logs/promtail-err.log