# build/supervisor/conf.d/demo.conf
[group:demo]
programs=demo-server

[program:demo-server]
directory=/app
command=sleep 3600
user=root
autostart=true
autorestart=true
stopsignal=TERM

stderr_logfile_maxbytes=1024MB
stdout_logfile_maxbytes=1024MB
stdout_logfile_backups=5
stdout_logfile=/data/logs/demo.log
stderr_logfile=/data/logs/demo-err.log